# etc/systemd/system/shared.automount
[Automount]
Where=/shared

[Install]
WantedBy=multi-user.target
